# s6/services.d/01-xplad/run
#!/usr/bin/with-contenv bash

/app/entrypoint.sh xplad start
